# resources/scripts/watercool.service
[Unit]
Description=Control the Kraken liquid cooler

[Service]
Type=simple
ExecStart=/usr/bin/watercool.py
